# deploy/stackobrite-agent.service
[Unit]
Description=Stackobrite Agent - Remote Execution Gateway
After=network.target containerd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/stackobrite-agent
Environment=AGENT_TOKEN_FILE=/etc/stackobrite/agent.key
Environment=AGENT_PORT=9090
Environment=AGENT_HOST=127.0.0.1
Environment=AGENT_KUBECONFIG=/etc/kubernetes/admin.conf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
